# bin/anthias-dev.service
[Unit]
Description=Anthias Dev Stack (Docker Compose)
Requires=docker.service network-online.target
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/pi/Anthias-Bellforge
ExecStart=/bin/bash -lc 'export MY_IP=$(hostname -I | awk '"'"'{print $1}'"'"'); /usr/bin/docker compose -f docker-compose.dev.yml -f docker-compose.viewer.yml up -d'
ExecStop=/usr/bin/docker compose -f docker-compose.dev.yml -f docker-compose.viewer.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
